# deploy/houseos.service
[Unit]
Description=HouseOS household operating system
After=local-fs.target

[Service]
Type=simple
User=@HOUSEOS_USER@
Group=@HOUSEOS_USER@
WorkingDirectory=/opt/houseos
EnvironmentFile=-/etc/houseos.env
ExecStart=/usr/bin/node /opt/houseos/server/index.js
Restart=always
RestartSec=3
NoNewPrivileges=false
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/lib/houseos

[Install]
WantedBy=multi-user.target
